# templates/mimir.service.erb
[Unit]
Description=<%= scope['mimir::service_name'] %>
After=network.target

[Service]
Type=simple
ExecStart=<%= scope['mimir::install_dir'] %>/mimir > /dev/null 2>&1
ExecStop=<%= scope['mimir::install_dir'] %>/mimir > /dev/null 2>&1
